# supervisord.conf
;-----------------------------------------------------------------------------------------------------------------------
; supervisor config file.
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/usr/appl/boogoo-process/hhr_admin/supervisord/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/usr/appl/boogoo-process/hhr_admin/supervisord/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available
; The sample program section below shows all possible program subsection values.
; Create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:hhr_hub]
command=gunicorn -k gevent -c hhr/flaskcelery/gun_hub.conf hhr.flaskcelery.tasksrouter:flask_app             ; the program (relative uses PATH, can take args)
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/usr/appl/boogoo-process          ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
user=root                    ; setuid to this UNIX account to run the program
stdout_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_hub.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_hub.log    ; stderr log path, NONE for none; default AUTO

[program:hhr_monitor]
command=pypy3 -m hhr.monitor 
process_name=%(program_name)s 
numprocs=1
directory=/usr/appl/boogoo-process
autostart=true
autorestart=unexpected
user=root
stdout_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_monitor.log
stderr_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_monitor.log

[program:hhr_filesrv]
command=gunicorn -k gevent -c gun_filesrv.conf server:fsrv_app
process_name=%(program_name)s 
numprocs=1
directory=/usr/appl/boogoo-process/hhr/filesrv
autostart=true
autorestart=unexpected
user=root
stdout_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_filesrv.log
stderr_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_filesrv.log

[program:hhr_worker]
command=celery -A hhr.celeryapp worker --loglevel=info -Q server01_1,server01_2,server01_3,SHARE_Q,broadcast_tasks --heartbeat-interval 10 -E -c 16
process_name=%(program_name)s 
numprocs=1
directory=/usr/appl/boogoo-process
autostart=true
autorestart=unexpected
user=root
stdout_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_worker.log
stderr_logfile=/usr/appl/boogoo-process/hhr_admin/log/hhr_worker.log
;-----------------------------------------------------------------------------------------------------------------------
